# deploy/systemd/road-user-violation.service
[Unit]
Description=Road User Platform Violation Detection
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=roaduser
WorkingDirectory=/opt/road-user-platform
EnvironmentFile=/etc/road-user-platform/server-common.env
ExecStart=/opt/road-user-platform/.venv/bin/python /opt/road-user-platform/src/violation_detection/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
