# deploy/systemd/road-user-forwarder.service
[Unit]
Description=Road User Platform MQTT Forwarder
After=network-online.target mosquitto.service road-user-backend.service
Wants=network-online.target

[Service]
Type=simple
User=roaduser
WorkingDirectory=/opt/road-user-platform
EnvironmentFile=/etc/road-user-platform/server-common.env
ExecStart=/opt/road-user-platform/.venv/bin/python /opt/road-user-platform/src/data_streaming/mqtt_forwarder.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
